# root/etc/services.d/sima/run
#!/bin/sh
/MPD_sima-0.14.1/vmpd-sima -c /MPD_sima-0.14.1/sima.conf
